# launchdaemons/com.mihomo.tun.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mihomo.tun</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mihomo</string>
    <string>-d</string>
    <string>/Users/yourname/.config/mihomo</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/yourname/.config/mihomo/mihomo.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yourname/.config/mihomo/mihomo.err</string>
</dict>
</plist>